# scripts-server/usurper-log-watcher.service
[Unit]
Description=Usurper Reborn Log Watcher
After=usurper-mud.service
Wants=usurper-mud.service

[Service]
Type=simple
# Loads DISCORD_WEBHOOK_URL. File lives outside this repo, 0640 root:usurper.
# Sample at scripts-server/log-watcher.env.example.
EnvironmentFile=/etc/usurper/log-watcher.env
ExecStart=/opt/usurper/log-watcher.sh
Restart=always
RestartSec=10
User=usurper
Group=usurper
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
